# .config/systemd/user/udiskie.service
[Unit]
Description=udiskie

[Service]
ExecStart=/usr/bin/udiskie -N

[Install]
WantedBy=default.target
