# systemd/.config/systemd/user/random-wallpaper.service
[Unit]
Description=Wallpaper Slideshow
Wants=hyprpaper.service graphical-session.target
After=hyprpaper.service graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/env random-wallpaper
Slice=background-graphical.slice
