# .config/systemd/user/ssh-agent.service.d/override.conf
[Service]
ExecStartPost=systemctl --user set-environment "SSH_AUTH_SOCK=%t/ssh-agent"
ExecStopPost=systemctl --user unset-environment SSH_AUTH_SOCK
